# launchd/com.propedge.grade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.propedge.grade</string>
  <key>ProgramArguments</key><array>
    <string>/usr/bin/python3</string>
    <string>REPO_PATH/scripts/grade_today.py</string>
  </array>
  <key>StartCalendarInterval</key><dict>
    <key>Hour</key><integer>5</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key><string>/tmp/propedge_grade.log</string>
  <key>StandardErrorPath</key><string>/tmp/propedge_grade_err.log</string>
  <key>WorkingDirectory</key><string>REPO_PATH</string>
</dict></plist>
